# priv/erlpmd.target
[Unit]
Description=EPMD service.
Requires=epmd-settle.service
After=epmd-settle.service

[Install]
Alias=epmd.target

